# scripts/clothy-db-backup.service
[Unit]
Description=Clothy DB backup sent to Telegram

[Service]
Type=oneshot
ExecStart=/opt/clothy/scripts/backup-db.sh
